# docker/gunicorn.supervisor.conf
[program:gunicorn]
directory=/app/www
command=/usr/local/bin/gunicorn -c /app/gunicorn_conf.py or-tools_api -b localhost:8000 --reload
stdout_logfile=/app/logs/gunicorn-stdout.log
stdout_logfile_maxbytes=1MB
stderr_logfile=/app/logs/gunicorn-stderr.log
stderr_logfile_maxbytes=1MB
autostart=true
autorestart=true
